# php-worker/supervisord.d/websockets.conf
[program:websockets]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/webinar.ws/artisan websockets:serve --port=8443
autostart=true
autorestart=true
numprocs=1
user=laradock
redirect_stderr=true
stderr_logfile=/var/log/webinar.ws.err.log
stdout_logfile=/var/log/webinar.ws.out.log